# crontab
0 * * * * ~/speed-monitor-run.sh > /proc/1/fd/1 2>/proc/1/fd/2
30 5 * * * ~/speed-monitor-run.sh upload > /proc/1/fd/1 2>/proc/1/fd/2
